# telegram_service/telegram_service.service
[Unit]
Description=Telegram AI Agent Bridge Service
After=network.target

[Service]
Type=simple
User=yasin
# Load secrets such as TELEGRAM_BOT_TOKEN from an ignored local file.
EnvironmentFile=-/home/yasin/git/derbent/telegram_service/.env
# Ensure node-based AI CLI tools are found when running under systemd.
Environment=PATH=/home/yasin/.nvm/versions/node/v22.22.2/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/yasin
WorkingDirectory=/home/yasin/git/derbent/telegram_service
ExecStart=/usr/bin/python3 /home/yasin/git/derbent/telegram_service/telegram_bot.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
